# supervisor/yii2-scheduler-queue.conf
[program:yii2-scheduler-queue-worker]
command=/path/to/project/folder/bin/yii queue-scheduler
autostart=true
autorestart=true
user=your_user
redirect_stderr=true